# all-in-one/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=java -jar %(ENV_BACKEND_JAR)s
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    SPRING_PROFILES_ACTIVE="local",
    DATABASE_URL="jdbc:postgresql://localhost:5432/%(ENV_POSTGRES_DB)s",
    DATABASE_USERNAME="%(ENV_POSTGRES_USER)s",
    DATABASE_PASSWORD="%(ENV_POSTGRES_PASSWORD)s",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    ENCRYPTION_SECRET="%(ENV_ENCRYPTION_SECRET)s",
    APP_CORS_ALLOWED_ORIGINS="*"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
